# docker/gazebo-supervisor.conf
[program:gazebo]
command=gosu ubuntu /opt/foxy-mac-sim/launch_gazebo_turtlebot3.sh
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
